# packages/arch/vpn-passthrough@.service
[Unit]
Description=VPN Passthrough Tunnel %i
After=vpn-passthrough.service
Requires=vpn-passthrough.service

[Service]
Type=notify
NotifyAccess=main
ExecStart=/usr/bin/vpn-passthrough start-tunnel %i
Restart=on-failure
RestartSec=5
Group=vpn-passthrough

[Install]
WantedBy=multi-user.target
